# conf/tornadotest_supervisor.ini
[program:tornadotest]
command=python /usr/local/x/tornadotest/src/tornadotest.py %(process_num)d
process_name=%(process_num)d
numprocs=10
numprocs_start = 10090
directory=/usr/local/x/tornadotest/src/
autostart=true
autorestart=true
stdout_logfile=/usr/local/x/tornadotest/log/output.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/usr/local/x/tornadotest/log/error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10


